# systemd/twc-neurio-sim.service
[Unit]
Description=Tesla Wall Connector Neurio simulator on Moxa UPort
After=systemd-modules-load.service systemd-udev-trigger.service

[Service]
Type=simple

# Load the Moxa UPort kernel driver.  On the tested Debian 13 host this module
# was installed from Moxa's Kernel 6.x driver package.
ExecStartPre=/sbin/modprobe mxuport

# Wait until all 16 Moxa serial devices exist.  Physical Moxa port 1 is
# /dev/ttyMXUSB0, physical port 16 is /dev/ttyMXUSB15.
ExecStartPre=/bin/sh -c 'for n in $(seq 1 45); do ok=1; for i in $(seq 0 15); do [ -e /dev/ttyMXUSB$i ] || ok=0; done; [ "$ok" = "1" ] && exit 0; sleep 1; done; exit 1'

# Critical Moxa-specific step: set each port to RS485 2-wire mode.  Without
# this, setserial reports Port: 0x0000 (RS232) and the Wall Connector traffic
# appears dead.
ExecStartPre=/bin/sh -c 'for i in $(seq 0 15); do /bin/setserial /dev/ttyMXUSB$i port 0x1; done'

# Fail fast if the first ports did not actually accept RS-485 2-wire mode.
ExecStartPre=/bin/sh -c 'setserial -g /dev/ttyMXUSB0 /dev/ttyMXUSB1 /dev/ttyMXUSB2 /dev/ttyMXUSB3 | grep -q "Port: 0x0001"'

ExecStart=/opt/twc-neurio-sim/twc_neurio_sim.py --quiet
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
